# script/install/crontab
0 1 * * * /bin/bash -l -c '/home/danbooru/scripts/backup_db'
0 2 * * * /bin/bash -l -c '/home/danbooru/scripts/backup_db_to_s3'
0 3 * * * /bin/bash -l -c '/home/danbooru/scripts/prune_backup_dbs'
