# script/install/related_tag_worker_init
### BEGIN INIT INFO
# Provides:          reportbooru_related_tag
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

NAME=reportbooru_related_tag
CHUID=danbooru
LOGFILE=/var/log/reportbooru/related_tag_worker.log
DESC="related tag worker"
PIDFILE=/var/run/reportbooru/related_tag_worker.pid
ROOT_DIR=/var/www/reportbooru/current
DAEMON=script/services/related_tag_worker.rb
DAEMON_CMD="bundle exec ruby $DAEMON --pidfile=$PIDFILE >> $LOGFILE 2>&1"

case "$1" in
    start)
  echo -n "Starting daemon: "$NAME
  start-stop-daemon --start --pidfile $PIDFILE --chuid $CHUID --chdir $ROOT_DIR --exec /bin/bash -- -l -c "$DAEMON_CMD"
  echo "."
  ;;

    stop)
  echo -n "Stopping daemon: "$NAME
  start-stop-daemon --stop --pidfile $PIDFILE --oknodo --chuid $CHUID --remove-pidfile
  echo "."
  ;;

    restart)
  echo -n "Restarting daemon: "$NAME
  start-stop-daemon --stop --retry 30 --pidfile $PIDFILE --oknodo --chuid $CHUID --remove-pidfile > $LOGFILE
  start-stop-daemon --start --quiet --pidfile $PIDFILE --exec /bin/bash --chuid $CHUID --chdir $ROOT_DIR -- -l -c "$DAEMON_CMD"
  echo "."
  ;;

    *)
  echo "Usage: "$1" {start|stop|restart}"
  exit 1

esac

exit 0
